# deploy/supervisor/runtime/rednote-opportunity.conf
[program:rednote-opportunity]
directory=/root/rednote-spider
command=/bin/bash -lc "source /root/rednote-spider/.venv/bin/activate && bash scripts/run_scheduled_loop.sh opportunity"
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stdout_logfile=/root/rednote-spider/logs/supervisor-opportunity.log
stderr_logfile=/root/rednote-spider/logs/supervisor-opportunity.err.log
environment=PYTHONPATH="src"
