# deploy/supervisor/runtime/rednote-discover.conf
[program:rednote-discover]
directory=/root/rednote-spider
command=/bin/bash -lc "source /root/rednote-spider/.venv/bin/activate && bash scripts/run_scheduled_loop.sh discover"
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stdout_logfile=/root/rednote-spider/logs/supervisor-discover.log
stderr_logfile=/root/rednote-spider/logs/supervisor-discover.err.log
environment=PYTHONPATH="src"
