# deploy/supervisor/runtime/rednote-login-qr.conf
[program:rednote-login-qr]
directory=/root/rednote-spider
command=/bin/bash -lc "source /root/rednote-spider/.venv/bin/activate && python scripts/init_schema.py && python scripts/run_login_qr_worker.py"
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stdout_logfile=/root/rednote-spider/logs/supervisor-login-qr.log
stderr_logfile=/root/rednote-spider/logs/supervisor-login-qr.err.log
environment=PYTHONPATH="src"
