# deploy/supervisor/supervisor-login-qr.conf
[program:rednote-login-qr]
directory=__PROJECT_DIR__
command=/bin/bash -lc "source __PROJECT_DIR__/.venv/bin/activate && python scripts/init_schema.py && python scripts/run_login_qr_worker.py"
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stdout_logfile=__PROJECT_DIR__/logs/supervisor-login-qr.log
stderr_logfile=__PROJECT_DIR__/logs/supervisor-login-qr.err.log
environment=PYTHONPATH="src"
