# deploy/supervisor/supervisor-ui.conf
[program:rednote-ui]
directory=__PROJECT_DIR__
command=/bin/bash -lc "bash scripts/run_ui_server.sh"
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stdout_logfile=__PROJECT_DIR__/logs/supervisor-ui.log
stderr_logfile=__PROJECT_DIR__/logs/supervisor-ui.err.log
environment=PYTHONPATH="src",STREAMLIT_SERVER_ADDRESS="127.0.0.1",STREAMLIT_SERVER_PORT="8501",STREAMLIT_SERVER_HEADLESS="true"
