# deploy/scheduler/systemd/rednote-discover.timer
[Unit]
Description=Run rednote discover cycle every 15 minutes

[Timer]
OnBootSec=3m
OnUnitActiveSec=15m
Unit=rednote-discover.service
Persistent=true

[Install]
WantedBy=timers.target
