# deploy/scheduler/systemd/rednote-opportunity.timer
[Unit]
Description=Run rednote opportunity cycle every 10 minutes

[Timer]
OnBootSec=5m
OnUnitActiveSec=10m
Unit=rednote-opportunity.service
Persistent=true

[Install]
WantedBy=timers.target
