# deploy/scheduler/supervisor-discover.conf
[program:rednote-discover]
directory=__PROJECT_DIR__
command=/bin/bash -lc "bash scripts/run_scheduled_loop.sh discover"
autostart=true
autorestart=true
startsecs=3
stopsignal=TERM
stdout_logfile=__PROJECT_DIR__/logs/supervisor-discover.log
stderr_logfile=__PROJECT_DIR__/logs/supervisor-discover.err.log
environment=PYTHONPATH="src",SCHED_DISCOVER_LOOP_INTERVAL_SECONDS="900"
